# recovery/etc/rc.local
#!/bin/sh

# Fix docker - if 'rooter' is detect it'll disable accept.
#iptables -I DOCKER-USER -j ACCEPT
iptables -P FORWARD ACCEPT
