# sidecar/etc/rc.local
#!/bin/sh

touch /tmp/rc-sidecar.local

# Fix docker - if it detects a 'router', it'll disable accept.
#iptables -I DOCKER-USER -j ACCEPT
iptables -P FORWARD ACCEPT

echo "rc.local started"
ifconfig -a

